# config/supervisor.d-supervisor.conf
[supervisord]
nodaemon=true

[program:consul]
command=/usr/local/bin/consul agent -data-dir /tmp/consul --join 172.17.0.3

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
